# django_hautomation_suite/django_hautomation_suite/distro_hasuite.conf
[program:hasuite]
command = python PROCESS_WORKING_DIRmanage.py run_gunicorn -w 1 -b 0.0.0.0:6000
directory = PROCESS_WORKING_DIR
user = PROCESS_USER
autostart = true
autorestart = true